# supervisor.conf
[program:pdf2png]
command=node /var/pdf2png/pdf2png.js
autostart=true
autorestart=true
stderr_logfile=/var/log/pdf2png.err.log
stdout_logfile=/var/log/pdf2png.out.log
